# deploy/roles/whoisit/templates/whoisit.timer
[Unit]
Description=%p

[Timer]
OnCalendar=hourly
Persistent=yes
Unit=%p.service

[Install]
WantedBy=multi-user.target
